# dcm4chee-2.14.2/dcm4jboss-sar/src/etc/bin/dcm4chee.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>dcm4chee</string>
	<key>WorkingDirectory</key>
	<string>DCM4CHEE_HOME/bin</string>
	
	<key>EnvironmentVariables</key>
	<dict>
	<key>JAVA_HOME</key>
	<string>/System/Library/Frameworks/JavaVM.framework/Home</string>
	</dict>
	
	<key>ProgramArguments</key>
	<array>
		<string>DCM4CHEE_HOME/bin/run.sh</string>
		<string>-b0.0.0.0</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
